# supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=redis-server

[program:fastchat]
command=java -jar -server -Xmx512m /root/fastchat.jar

